# deploy/systemd/bybit-sync-docker.timer
[Unit]
Description=Daily Bybit Data Sync (Docker)

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
